# ansible/roles/monitoring/templates/backup.service.j2
[Unit]
Description=Бэкап данных Docker Monitoring

[Service]
Type=oneshot
ExecStart=/srv/docker/monitoring/backup.sh
